# dashboard/deploy/crossfit-metropolitano-dashboard.service
[Unit]
Description=CrossFit MetroPolitano churn dashboard
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/crossfit-metropolitano-dashboard
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/crossfit-metropolitano-dashboard/dashboard_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

